# source/minibroker-aprs.service
[Unit]
Description=MiniBroker APRS Bridge. Version 2.1
After=network.target minibroker-emergencias.service
Wants=minibroker-emergencias.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=meshnet
Group=meshnet
WorkingDirectory=/opt/minibroker
EnvironmentFile=/opt/minibroker/.env-usb
ExecStart=/opt/minibroker/venv/bin/python /opt/minibroker/emergency_aprs_bridge.py

# El bridge debe seguir vivo y reconectar por sí mismo.
# No debe reiniciarse por una parada "limpia" inducida por dependencias externas.
Restart=on-failure
RestartSec=3


[Install]
WantedBy=multi-user.target